# conf/supervisord.conf
[supervisord]
nodaemon=true

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:uwsgi]
command = uwsgi --ini /etc/uwsgi.ini
directory=/opt/jumper
autostart=true

[program:nginx]
command = nginx -g "daemon off;"
autostart=true
